# mqtt-to-prom/meta/mqtt-to-prom.service
[Unit]
Description=mqtt to prometheus
After=network.target

[Service]
ExecStart=/usr/bin/mqtt-to-prom

[Install]
WantedBy=multi-user.target
